# web/deploy/medicart-tunnel.service
[Unit]
Description=cloudflared tunnel — intel2.thatshoon.com (MediCart)
After=network.target medicart-backend.service medicart-frontend.service

[Service]
User=rokey
ExecStart=/home/rokey/.local/bin/cloudflared --no-autoupdate --config /home/rokey/.cloudflared/config.yml tunnel run intel2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
